# systemd/phlex-server.service
[Unit]
Description=Phlex Media Server
Documentation=https://docs.phlex.media
After=network.target mysql.service
Wants=mysql.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
User=phlex
Group=phlex
WorkingDirectory=/var/www/phlex
ExecStart=/usr/bin/php /var/www/phlex/public/index.php start
ExecReload=/bin/kill -SIGUSR1 $MAINPID
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30
TimeoutStartSec=30

Environment="PHLEX_ENV=production"
EnvironmentFile=/etc/phlex/env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=phlex-server

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/phlex /var/log/phlex /run
RestrictNamespaces=true
LockPersonality=true
RemoveIPC=true

[Install]
WantedBy=multi-user.target
